# scripts/systemd/hexevoice-stt.service.in
[Unit]
Description=HexeVoice faster-whisper STT

[Service]
WorkingDirectory=__ROOT_DIR__
LimitNOFILE=65536
ExecStart=/bin/bash -lc 'source "__ENV_FILE__"; cd "__ROOT_DIR__"; eval "$STT_CMD"'
