# scripts/systemd/hexevoice-backend.service.in
[Unit]
Description=HexeVoice Backend

[Service]
WorkingDirectory=__ROOT_DIR__
ExecStart=/bin/bash -lc 'source "__ENV_FILE__"; exec __BACKEND_CMD__'
